# docker/images/openam/supervisor-apache.conf
[program:tomcat]
command=service apache2 start
user=root
autostart=false
autorestart=unexpected
