# docker/images/openam/supervisor-tomcat.conf
[program:tomcat]
command=service tomcat8 start
user=root
autostart=true
autorestart=unexpected
