# docker/images/openam/supervisor-tomcat.conf
[program:tomcat]
command=/bin/bash -c "service tomcat8 start"
user=root
autostart=true
autorestart=false
