# docker/images/promed-scraper/crontab
0 0 * * * bash /cronjob.sh >> /cron.log 2>&1
