# docker/images/promed-scraper/crontab
0 0 * * * /cronscript.sh >> /cron.log 2>&1
