# docker/images/mobile-survey/supervisor-mobile-survey.conf
[program:mobile-survey]
command=/bin/bash /mobile-survey.sh
user=root
autostart=true
autorestart=unexpected
