# docker/images/spa/spa-supervisor.conf
[program:spa]
command=/bin/bash /spa.sh
user=root
autostart=true
autorestart=unexpected
