# docker/images/spa/spa-supervisor.conf
[program:flirt]
command=/bin/bash /flirt.sh
user=root
autostart=true
autorestart=unexpected

