# docker/images/niam/supervisord.conf
[program:niam]
command=meteor run
user=root
autostart=true
autorestart=unexpected
directory=/nia-monitor
