# docker/images/niam/supervisord.conf
[program:niam]
command=/bin/bash /niam.sh
user=root
autostart=true
autorestart=unexpected
