# docker/images/era/era-supervisor.conf
[program:era]
command=/bin/bash /era.sh
user=root
autostart=true
autorestart=unexpected
