# docker/images/grits/diagnostic-dashboard/dashboardd.conf
[program:dashboard]
command=/dashboard.sh
user=grits
numprocs=1
autostart=true
autorestart=true
startsecs=10
logfile_maxbytes=150MB
logfile_backups=10
loglevel=info

