# docker/images/grits/diagnostic-dashboard/dashboardd.conf
[program:dashboard]
command=node %(GRITS_HOME)s/diagnostic-dashboard/bundle/main.js
environment=METEOR_DB_NAME=%(meteor_db_name)s,MONGO_URL=%(MONGO_URL)s,PORT=%(meteor_port)s,ROOT_URL=%(root_url)s,MAIL_URL=%(mail_url)s",GRID_URL="%(grid_url)s",BSVE_ACCESS_KEY=%(bsve_access_key)s
directory=%(GRITS_HOME)s/diagnostic-dashboard
user=grits
numprocs=1
autostart=true
autorestart=true
startsecs=10
logfile_maxbytes=150MB
logfile_backups=10
loglevel=info

