# docker/images/grits/diagnostic-dashboard/dashboardd.conf
[program:dashboard]
command=/dashboard.sh
user=root
numprocs=1
autostart=%(ENV_AUTOSTART_DASHBOARD)s
autorestart=true
startsecs=10
logfile_maxbytes=150MB
logfile_backups=10
loglevel=info

