# docker/images/flirt/celery-supervisor.conf
[program:celery]
command=/bin/bash /celery.sh
user=root
autostart=true
autorestart=unexpected

