# docker/images/flirt/flirt-supervisor.conf
[program:flirt]
command=node /example/gritsbuild/bundle/main.js 
user=root
autostart=true
autorestart=unexpected

