# docker/images/flirt/flirt-supervisor.conf
[program:flirt]
command=METEOR_SETTINGS=$(cat /shared/settings-production.json) node /example/gritsbuild/bundle/main.js 
user=root
autostart=true
autorestart=unexpected

