# docker/images/flirt/tornado-supervisor.conf
[program:tornado]
command=/bin/bash /tornado.sh
user=root
autostart=true
autorestart=unexpected

