# docker/images/grits/grits-api/gritsapid.conf
; We use a group like this so that the grits-api deploy
; script doesn't need to know the names of all the supervisor programs.
[group:gritsapigroup]
programs=celery_priority,celery_process,celery_diagnose,gritsapi
priority=999

[program:gritsapi]
command=%(GRITS_HOME)s/grits_env/bin/python server.py
directory=%(GRITS_HOME)s/grits-api
user=grits
environment=NLTK_DATA=%(GRITS_HOME)s/nltk_data,MONGO_URL=%(MONGO_URL)s
numprocs=1
autostart=true
autorestart=true
startsecs=10
logfile_maxbytes=150MB
logfile_backups=10
loglevel=info
