# docker/images/parse-server/supervisor-dashboard.conf
[program:dashboard]
command=/bin/bash /dashboard.sh
user=root
autostart=true
autorestart=unexpected
