# docker/images/parse-server/supervisor-server.conf
[program:server]
command=/bin/bash /parse-server.sh
user=root
autostart=true
autorestart=unexpected
