# docker/images/grits/logging/supervisor-logging.conf
[program:logging]
command=/logging.sh
user=root
numprocs=1
autostart=false
autorestart=true
startsecs=10
logfile_maxbytes=1MB
loglevel=info


