# docker/images/app-router/supervisor-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -f /etc/ssh/sshd_config -D
autostart=true
autorestart=unexpected
