# docker/images/app-router/supervisor-haproxyd.conf
[program:haproxyd]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg -d 
autostart=true
autorestart=unexpected
