# docker/images/birt/birt-supervisor.conf
[program:birt]
command=/bin/bash /birt.sh
user=root
autostart=true
autorestart=unexpected
